# cron.d/stock
10 18  * * 1-5 root test -x /var/stock/run.sh && /var/stock/run.sh
15 09  * * 2-6 root test -x /var/stock/update_N225.sh && /var/stock/update_N225.sh
