# cron.d/stock
10 18  * * 1-5 root test -x /var/stock/run.sh && /var/stock/run.sh
